# docker/images/water/supervisor.d/supervisor.ini
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
logfile_backups=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off; error_log /dev/stderr info;"
autostart=true
environment = PYTHONUNBUFFERED=1
autorestart=true
priority=10
redirect_stderr=true
redirect_stdout=true
stdout_logfile=NONE
stderr_logfile=NONE
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
stderr_logfile_backups=0
stdout_logfile_backups=0
stdout_events_enabled = true
stderr_events_enabled = true

[program:php-fpm]
command=/usr/bin/php-fpm -c /etc/php/php.ini --fpm-config /etc/php/php-fpm.conf -F -O
numprocs=1
environment = PYTHONUNBUFFERED=1
user=root
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true
redirect_stdout=true
stdout_logfile=NONE
stderr_logfile=NONE
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
stderr_logfile_backups=0
stdout_logfile_backups=0
stdout_events_enabled = true
stderr_events_enabled = true

[program:email-queue]
command=php /var/www/trustwaiver/artisan queue:work --sleep=3 --tries=3
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/trustwaiver/storage/logs/supervisord-email-queue.log
